# etc/raspberrypios/wrolpi-repair.service
[Unit]
Description=WROLPi Repair Service
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/wrolpi/repair.sh
User=root
StandardOutput=journal
StandardError=journal
