# debian-live-config/config/includes.chroot/etc/xdg/autostart/wrolpi-firstboot-notify.desktop
[Desktop Entry]
Type=Application
Name=WROLPi First-Boot Notification
Comment=Show a summary of what WROLPi bootstrap did on first boot
Exec=/usr/local/bin/wrolpi-firstboot-notify
StartupNotify=false
NoDisplay=true
