# debian-live-config/config/includes.chroot/etc/xdg/autostart/wrolpi-firstboot-notify.desktop
[Desktop Entry]
Type=Application
Name=WROLPi First-Boot Notification
Comment=Show a summary of what WROLPi bootstrap did on first boot
Exec=/usr/local/bin/wrolpi-firstboot-notify
TryExec=/usr/local/bin/wrolpi-firstboot-notify
Icon=dialog-information
Terminal=false
Categories=System;
NoDisplay=true
X-GNOME-Autostart-enabled=true
