# etc/raspberrypios/postgresql@15-map.service.d/override.conf
[Service]
# WROLPi's map uses a lot of processing power.  Change niceness of map database to prevent blocking of other services.
Nice=19